# packaging/alpine_3.6/jobber.post-deinstall
#!/bin/sh

# delete client user
deluser jobber_client

# remove service
rc-update del jobber
