# packaging/alpine_3.6/jobber.post-install
#!/bin/sh

DAEMON=/usr/libexec/jobbermaster

# set owner and perms of binaries
chown root:root "${DAEMON}"
chmod 0755 "${DAEMON}"

# start jobber
rc-update add jobber
service jobber start
