# packaging/alpine/jobber.post-deinstall
#!/bin/sh

# remove service
rc-update del jobber
